# deploy/oracle/caddy-proxy.service
[Unit]
Description=Caddy HTTPS reverse proxy for Sales Intelligence
After=network-online.target sales-intelligence.service
Wants=network-online.target

[Service]
Type=simple
Environment=HOME=/var/snap/caddy/common
WorkingDirectory=/var/snap/caddy/common
ExecStart=/snap/caddy/current/usr/bin/caddy run --config /var/snap/caddy/common/Caddyfile --adapter caddyfile
ExecReload=/snap/caddy/current/usr/bin/caddy reload --config /var/snap/caddy/common/Caddyfile --adapter caddyfile --force
TimeoutStopSec=5s
LimitNOFILE=1048576
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
