# deploy/oracle/sales-intelligence.service
[Unit]
Description=Sales Intelligence Streamlit app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/app
Environment=PATH=/home/ubuntu/app/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/app/venv/bin/streamlit run app.py --server.address 0.0.0.0 --server.port 8501 --browser.gatherUsageStats false
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
